# user_config/systemd/user/sunshine.service.d/env.conf
[Service]
Environment=DISPLAY=:0
#Environment=DISPLAY=:101
#Environment=XAUTHORITY=%h/.Xauthority
#Environment=XDG_RUNTIME_DIR=/run/user/%U

